# ansible/02-ClusterInstall/roles/iperf/templates/iperf.service.j2
[Unit]
Description=iperf - network measurement tool
After=network.target

[Service]
ExecStart=/bin/iperf3 -s

[Install]
WantedBy=multi-user.target